# docker/etc/supervisor.d/app-server.ini
[program:app-server]
process_name=app-server
autostart=true
autorestart=true
startsecs=10
user=root
directory=/home/code/api
command=node index.js
stderr_logfile = /var/log/supervisord/err.log
stdout_logfile = /var/log/supervisord/out.log